# python/homelab-monitor/v2-systemd/homelab-monitor.timer
[Unit]
Description=Run the homelab monitor every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=1min
Unit=homelab-monitor.service

[Install]
WantedBy=timers.target
